# deploy/systemd/release-webhook-receiver.service
[Unit]
Description=SLO Rollout Release Webhook Receiver
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/slo-rollout-demo
Environment=OLLAMA_URL=http://192.168.30.1:11434
Environment=MODEL=qwen2.5:3b
ExecStart=/usr/bin/python3 /root/slo-rollout-demo/scripts/release-webhook-receiver.py --host 0.0.0.0 --port 18080
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
