# deploy/systemd/release-rollout-watcher.service
[Unit]
Description=SLO Rollout Go Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/slo-rollout-demo
ExecStart=/root/slo-rollout-demo/watcher/release-rollout-watcher \
  --namespace slo-rollout \
  --rollout demo-app \
  --repo-dir /root/slo-rollout-demo \
  --ollama-url http://192.168.30.1:11434 \
  --model qwen2.5:3b \
  --interval 10s
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
